# project/main.cron
10 */1 * * *	/usr/bin/python /home/de1h/S_00/main.py								# запуск программы main.py
0 5 * * *		find /home/de1h/S_00/*.bz2 -mtime +90 -delete						# удаление архивных файлов
0 7 * * 7		find /home/de1h/S_00/log/*log -mtime +7 -exec bzip2 {} \;			# архивирование журналов (логи)